# res/inittab
::sysinit:/etc/init.d/rcS
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::askfirst:-/bin/sh
tty2::respawn:/bin/sh
tty3::respawn:/bin/sh
tty4::respawn:/bin/sh
tty5::respawn:/bin/sh
